# DEBIAN/postinst
#!/bin/sh
set -e
if [ -x "/etc/init.d/serclient" ]; then
	update-rc.d serclient defaults >/dev/null
	invoke-rc.d serclient start || exit $?
fi
