# DEBIAN/prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/serclient" ]; then
	invoke-rc.d serclient stop || exit $?
fi
